# systemd/utils/minecraft-restart@.service
[Unit]
Description=Restart Minecraft server %i

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-restart minecraft@%i.service
